# init/corosync-qnetd.service.in
[Unit]
Description=Corosync Qdevice Network daemon
Documentation=man:corosync-qnetd
ConditionKernelCommandLine=!nocluster
Requires=network-online.target
After=network-online.target

[Service]
EnvironmentFile=@SYSCONFDIR@/sysconfig/corosync-qnetd
ExecStart=@BINDIR@/corosync-qnetd -f $COROSYNC_QNETD_OPTIONS
Type=simple
Restart=on-abnormal
# Uncomment and set user who should be used for executing qnetd
#User=coroqnetd

[Install]
WantedBy=multi-user.target
