# init/corosync-notifyd.service.in
[Unit]
Description=Corosync Dbus and snmp notifier
Documentation=man:corosync-notifyd
Wants=corosync.service
After=corosync.service

[Service]
EnvironmentFile=@SYSCONFDIR@/sysconfig/corosync-notifyd
ExecStart=@SBINDIR@/corosync-notifyd -f $OPTIONS
Type=simple
Restart=on-failure

[Install]
WantedBy=multi-user.target
