# init/corosync-qdevice.service.in
[Unit]
Description=Corosync Qdevice daemon
Documentation=man:corosync-qdevice
ConditionKernelCommandLine=!nocluster
Wants=corosync.service
After=corosync.service

[Service]
EnvironmentFile=-@INITCONFIGDIR@/corosync-qdevice
ExecStart=@SBINDIR@/corosync-qdevice -f $COROSYNC_QDEVICE_OPTIONS
Type=notify
RuntimeDirectory=corosync-qdevice
RuntimeDirectoryMode=0770
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
