# init/corosync-qnetd.service.in
[Unit]
Description=Corosync Qdevice Network daemon
ConditionKernelCommandLine=!nocluster
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=@INITWRAPPERSDIR@/corosync-qnetd start
ExecStop=@INITWRAPPERSDIR@/corosync-qnetd stop
Type=forking

[Install]
WantedBy=multi-user.target
