# init/corosync-notifyd.service.in
[Unit]
Description=Corosync Dbus and snmp notifier
Requires=corosync.service
After=corosync.service

[Service]
EnvironmentFile=-@SYSCONFDIR@/sysconfig/corosync-notifyd
EnvironmentFile=-@SYSCONFDIR@/default/corosync-notifyd
ExecStart=@SBINDIR@/corosync-notifyd -f $OPTIONS
Type=simple

[Install]
WantedBy=multi-user.target
