# init/corosync-qdevice.service.in
[Unit]
Description=Corosync Qdevice daemon
Documentation=man:corosync-qdevice
ConditionKernelCommandLine=!nocluster
Wants=corosync.service
After=corosync.service

[Service]
ExecStart=@INITWRAPPERSDIR@/corosync-qdevice start
ExecStop=@INITWRAPPERSDIR@/corosync-qdevice stop
Type=forking

[Install]
WantedBy=multi-user.target
