# init/openrc/corosync-notifyd.in
#!/sbin/openrc-run

pidfile="/run/$RC_SVCNAME.pid"
command="/usr/sbin/corosync-notifyd"

depend() {
        need net
}
